# maintenance/cron.d/cron
0 0 * * * echo "Renewed certificates: `date`" > /var/log/cron
0 0 * * * certbot renew
